# unified_deployment/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/nginx/stdout.log
stderr_logfile=/var/log/nginx/stderr.log

# Landing Page - Port 8501
[program:landing]
command=streamlit run /app/landing_page.py --server.port=8501 --server.address=0.0.0.0 --server.headless=true --server.baseUrlPath=/
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/landing.log
stderr_logfile=/var/log/supervisor/landing_err.log

# Access Controller Agent - Port 8502
[program:access_controller]
command=streamlit run streamlit_app.py --server.port=8502 --server.address=0.0.0.0 --server.headless=true --server.baseUrlPath=/access-controller
directory=/app/access_controller_agent
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/access_controller.log
stderr_logfile=/var/log/supervisor/access_controller_err.log

# Portfolio Manager - Port 8503
[program:portfolio_manager]
command=streamlit run streamlit_app.py --server.port=8503 --server.address=0.0.0.0 --server.headless=true --server.baseUrlPath=/portfolio-manager
directory=/app/portfolio_manager
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/portfolio_manager.log
stderr_logfile=/var/log/supervisor/portfolio_manager_err.log

# Campaign Validator - Port 8504
[program:campaign_validator]
command=streamlit run app.py --server.port=8504 --server.address=0.0.0.0 --server.headless=true --server.baseUrlPath=/campaign-validator
directory=/app/campaign_validator_agent
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/campaign_validator.log
stderr_logfile=/var/log/supervisor/campaign_validator_err.log

# PR Code Reviewer - Port 8505
[program:pr_code_reviewer]
command=streamlit run streamlit_app.py --server.port=8505 --server.address=0.0.0.0 --server.headless=true --server.baseUrlPath=/pr-code-reviewer
directory=/app/pr_code_reviewer
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/pr_code_reviewer.log
stderr_logfile=/var/log/supervisor/pr_code_reviewer_err.log

# Research Assistant - Port 8506
[program:research_assistant]
command=streamlit run streamlit_app.py --server.port=8506 --server.address=0.0.0.0 --server.headless=true --server.baseUrlPath=/research-assistant
directory=/app/research_assistant
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/research_assistant.log
stderr_logfile=/var/log/supervisor/research_assistant_err.log

# Loan Underwriter - Port 8507
[program:loan_underwriter]
command=streamlit run streamlit_app.py --server.port=8507 --server.address=0.0.0.0 --server.headless=true --server.baseUrlPath=/loan-underwriter
directory=/app/loan_underwriter
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/loan_underwriter.log
stderr_logfile=/var/log/supervisor/loan_underwriter_err.log

# Content Moderator - Port 8508
[program:content_moderator]
command=streamlit run streamlit_app.py --server.port=8508 --server.address=0.0.0.0 --server.headless=true --server.baseUrlPath=/content-moderator
directory=/app/content_moderator
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/content_moderator.log
stderr_logfile=/var/log/supervisor/content_moderator_err.log

# HelpDesk Bot - Port 8509
[program:helpdesk_bot]
command=streamlit run streamlit_app.py --server.port=8509 --server.address=0.0.0.0 --server.headless=true --server.baseUrlPath=/helpdesk-bot
directory=/app/helpdesk_bot
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/helpdesk_bot.log
stderr_logfile=/var/log/supervisor/helpdesk_bot_err.log

# Social Media Trends - Port 8510
[program:social_media_trends]
command=streamlit run streamlit_app.py --server.port=8510 --server.address=0.0.0.0 --server.headless=true --server.baseUrlPath=/social-media-trends
directory=/app/social_media_trends
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/social_media_trends.log
stderr_logfile=/var/log/supervisor/social_media_trends_err.log

# Portfolio Swarm - Port 8511
[program:portfolio_swarm]
command=streamlit run streamlit_app.py --server.port=8511 --server.address=0.0.0.0 --server.headless=true --server.baseUrlPath=/portfolio-swarm
directory=/app/portfolio_swarm
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/portfolio_swarm.log
stderr_logfile=/var/log/supervisor/portfolio_swarm_err.log
